# publish/Data/crontab
*/20 *   *   *   *    curl -s "http://127.0.0.1:9117/jsondb/save"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/torrentby/parse"
1    */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/torrentby/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9117/cron/torrentby/ParseAllTask"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/rutracker/parse"
1    */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/rutracker/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9117/cron/rutracker/ParseAllTask"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/rutor/parse"
5    */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/rutor/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9117/cron/rutor/ParseAllTask"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/megapeer/parse"
5    */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/megapeer/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9117/cron/megapeer/ParseAllTask"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/nnmclub/parse"
10   */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/nnmclub/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9117/cron/nnmclub/ParseAllTask"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/kinozal/parse"
10   */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/kinozal/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9117/cron/kinozal/ParseAllTask"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/bitru/parse"
20   */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/bitru/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9117/cron/bitru/ParseAllTask"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/selezen/parse"
20   */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/selezen/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9117/cron/selezen/ParseAllTask"

*/40 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/toloka/parse"
20   */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/toloka/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9117/cron/toloka/ParseAllTask"

40   */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/anilibria/parse?limit=1400"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/animelayer/parse"
50   */8 *   *   *    curl -s "http://127.0.0.1:9117/cron/animelayer/parse?maxpage=90"

*/40 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/anifilm/parse"
#50   1   *   *   *    curl -s "http://127.0.0.1:9117/cron/anifilm/parse?fullparse=true"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/hdrezka/parse"
#50   1   *   *   *    curl -s "http://127.0.0.1:9117/cron/hdrezka/parse?maxpage=83"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/baibako/parse"
#50   1   *   *   *    curl -s "http://127.0.0.1:9117/cron/baibako/parse?maxpage=418"

*/40 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/lostfilm/parse"
#1    1   */5 *   *    curl -s "http://127.0.0.1:9117/cron/lostfilm/parse?maxpage=1000"
